# deployment/cadre-data.supervisor.ini
[program:cadre-data]
command = /var/www/cadre-data/start.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/cadre-data.super.err.log
stdout_logfile=/var/log/cadre-data.super.out.log